# supervisor/supervisord-newrelic-plugin-agent.conf
[program:newrelic-plugin-agent]
command=/usr/local/bin/newrelic-plugin-agent -f -c "/etc/newrelic/newrelic-plugin-agent.cfg"
numprocs=1
autostart=false
autorestart=true
