# supervisor/supervisord-newrelic-plugin-agent-discover.conf
[program:newrelic-plugin-agent-configuration]
command=/etc/newrelic/newrelic-plugin-agent-discover.sh
numprocs=1
autostart=true
autorestart=true
